# monitoring-server/platforms/healthyhome-rpi/monitd.service
[Unit]
Description=HealthyHome monitoring server

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/monitd

[Install]
WantedBy=multi-user.target